# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "Execute cron service."
exec crond -f
